# timelapsepi.service
[Unit]
Description=TimelapsePI Camera Controller
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/timelapsepi
ExecStart=/usr/bin/python3 /home/pi/timelapsepi/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
